# backend/sigma_rules/builtin_sigma.json
[
  {
    "id": "sigma-001",
    "name": "Suspicious PowerShell Encoded Command",
    "description": "Detects PowerShell with encoded command and evasion flags",
    "severity": "high",
    "mitre_technique": "T1059.001",
    "mitre_tactic": "Execution",
    "source": "event_log",
    "event_ids": [4688, 4104],
    "patterns": ["powershell.*-enc", "powershell.*-encodedcommand", "powershell.*-nop.*-w.*hidden"]
  },
  {
    "id": "sigma-002",
    "name": "LSASS Memory Access",
    "description": "Detects access to LSASS process memory for credential dumping",
    "severity": "critical",
    "mitre_technique": "T1003.001",
    "mitre_tactic": "Credential Access",
    "source": "event_log",
    "event_ids": [4656, 4663, 10],
    "patterns": ["lsass.exe", "lsass.dmp"]
  },
  {
    "id": "sigma-003",
    "name": "Scheduled Task Created via CMD",
    "description": "Detects scheduled task creation using schtasks",
    "severity": "high",
    "mitre_technique": "T1053.005",
    "mitre_tactic": "Persistence",
    "source": "event_log",
    "event_ids": [4698, 4702],
    "patterns": ["schtasks.*create", "schtasks.*change"]
  },
  {
    "id": "sigma-004",
    "name": "New Service Installed",
    "description": "Detects installation of a new Windows service",
    "severity": "high",
    "mitre_technique": "T1543.003",
    "mitre_tactic": "Persistence",
    "source": "event_log",
    "event_ids": [7045, 4697],
    "patterns": ["sc.*create", "New-Service"]
  },
  {
    "id": "sigma-005",
    "name": "Suspicious MSHTA Execution",
    "description": "Detects mshta.exe executing remote payloads",
    "severity": "critical",
    "mitre_technique": "T1218.005",
    "mitre_tactic": "Defense Evasion",
    "source": "event_log",
    "event_ids": [4688],
    "patterns": ["mshta.exe.*http", "mshta.*vbscript", "mshta.*javascript"]
  },
  {
    "id": "sigma-006",
    "name": "WMI Remote Command Execution",
    "description": "Detects WMI used for remote process execution",
    "severity": "high",
    "mitre_technique": "T1047",
    "mitre_tactic": "Execution",
    "source": "event_log",
    "event_ids": [4688],
    "patterns": ["wmic.*process.*call.*create", "wmic.*/node:"]
  },
  {
    "id": "sigma-007",
    "name": "CertUtil Used for Download",
    "description": "Detects certutil.exe used to download files (LOLBin)",
    "severity": "high",
    "mitre_technique": "T1105",
    "mitre_tactic": "Command and Control",
    "source": "event_log",
    "event_ids": [4688],
    "patterns": ["certutil.*urlcache", "certutil.*-decode", "certutil.*http"]
  },
  {
    "id": "sigma-008",
    "name": "Registry SAM Hive Export",
    "description": "Detects attempts to export SAM/SYSTEM registry hives",
    "severity": "critical",
    "mitre_technique": "T1003.002",
    "mitre_tactic": "Credential Access",
    "source": "event_log",
    "event_ids": [4663, 4656],
    "patterns": ["reg.*export.*sam", "reg.*save.*security", "reg.*save.*system"]
  },
  {
    "id": "sigma-009",
    "name": "Suspicious Regsvr32 Execution",
    "description": "Detects regsvr32 used to load remote DLLs",
    "severity": "high",
    "mitre_technique": "T1218.010",
    "mitre_tactic": "Defense Evasion",
    "source": "event_log",
    "event_ids": [4688],
    "patterns": ["regsvr32.*/s.*/n.*/u.*/i:http", "regsvr32.*scrobj"]
  },
  {
    "id": "sigma-010",
    "name": "Event Log Cleared",
    "description": "Detects Windows event log being cleared (anti-forensics)",
    "severity": "critical",
    "mitre_technique": "T1070.001",
    "mitre_tactic": "Defense Evasion",
    "source": "event_log",
    "event_ids": [1102, 104],
    "patterns": ["Log cleared", "Security Log", "System Log"]
  },
  {
    "id": "sigma-011",
    "name": "New Local Admin Account Created",
    "description": "Detects creation of new local administrator accounts",
    "severity": "high",
    "mitre_technique": "T1136.001",
    "mitre_tactic": "Persistence",
    "source": "event_log",
    "event_ids": [4720, 4732],
    "patterns": ["net user.*add", "net localgroup.*administrators.*add"]
  },
  {
    "id": "sigma-012",
    "name": "Pass-The-Hash Activity",
    "description": "Detects Pass-The-Hash lateral movement pattern",
    "severity": "critical",
    "mitre_technique": "T1550.002",
    "mitre_tactic": "Lateral Movement",
    "source": "event_log",
    "event_ids": [4624, 4648, 4776],
    "patterns": ["NTLM.*Logon.*Type.*3", "Anonymous.*NTLM"]
  }
]
